# examples/fargate-envoy-spire/cli-client-with-envoy/supervisord.conf
[supervisord]
nodaemon=true
loglevel=info
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisord.pid

[program:spire-agent]
command=/usr/local/bin/spire-agent run -config /opt/spire/conf/agent/agent.conf
directory=/opt/spire
autostart=true
autorestart=true
stderr_logfile=/var/log/supervisor/spire-agent.err.log
stdout_logfile=/var/log/supervisor/spire-agent.out.log
user=root
priority=100

[program:envoy]
command=/usr/local/bin/envoy -l info -c /etc/envoy/envoy.yaml
directory=/etc/envoy
autostart=true
autorestart=true
stderr_logfile=/var/log/supervisor/envoy.err.log
stdout_logfile=/var/log/supervisor/envoy.out.log
user=root
priority=200
# Start Envoy after SPIRE agent is running
depends_on=spire-agent

# Note: CLI client runs on-demand, not as a persistent service
# It can be executed manually using: docker exec <container> /usr/local/bin/spire-cli [options]
