# examples/fargate-envoy-spire/api-service-with-envoy/supervisord.conf
[supervisord]
nodaemon=true
loglevel=info
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisord.pid

[program:spire-agent]
command=/usr/local/bin/spire-agent run -config /opt/spire/conf/agent/agent.conf
directory=/opt/spire
autostart=true
autorestart=true
stderr_logfile=/var/log/supervisor/spire-agent.err.log
stdout_logfile=/var/log/supervisor/spire-agent.out.log
user=root
priority=100

[program:envoy]
command=/usr/local/bin/envoy -l info -c /etc/envoy/envoy.yaml
directory=/etc/envoy
autostart=true
autorestart=true
stderr_logfile=/var/log/supervisor/envoy.err.log
stdout_logfile=/var/log/supervisor/envoy.out.log
user=root
priority=200
# Start Envoy after SPIRE agent is running
depends_on=spire-agent

[program:api-service]
command=/usr/local/bin/api-service
directory=/usr/local/bin
autostart=true
autorestart=true
stderr_logfile=/var/log/supervisor/api-service.err.log
stdout_logfile=/var/log/supervisor/api-service.out.log
environment=PORT="8082"
user=root
priority=300
# Start API service after both SPIRE agent and Envoy are running
depends_on=spire-agent,envoy
